# deps/emerge_sync_cron/crontab
# update portage every 12 hours
* 6,18 * * * /usr/bin/emerge --sync
